# deepamehta-4.0.8/debian/postinst
#!/bin/sh
# postinst script for deepamehta
#
# see: dh_installdeb(1)

set -e
#. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    chmod 755 /etc/init.d/deepamehta
    chmod 755 /usr/lib/deepamehta/deepamehta-daemon
    chmod 644 /etc/deepamehta/*
    chmod 644 /etc/default/deepamehta
    update-rc.d deepamehta defaults
    chown -R deepamehta:deepamehta /var/lib/deepamehta
    chown -R deepamehta:deepamehta /var/cache/deepamehta
    chown -R deepamehta:deepamehta /var/log/deepamehta
    chown -R deepamehta:deepamehta /var/run/deepamehta
    cp /usr/share/doc/deepamehta/examples/deepamehta.filedir /var/lib/deepamehta/deepamehta-filedir/README.txt
    /etc/init.d/deepamehta start
    ;;
    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
